# deploy/supervisord.conf
; Sample supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Notes:
;  - Shell expansion ("~" or "$HOME") is not supported.  Environment
;    variables can be expanded using this syntax: "%(ENV_HOME)s".
;  - Comments must have a leading space: "a=b ;comment" not "a=b;comment".

[unix_http_server]
file=/home/jabberbot/supervisor.sock   ; (the path to the socket file)

[inet_http_server]
port=127.0.0.1:9004

[supervisord]
logfile=/home/jabberbot/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/home/jabberbot/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///home/jabberbot/supervisor.sock ; use a unix:// URL  for a unix socket

;[include]
;files = relative/directory/*.ini

[program:bot]
command=/home/jabberbot/botenv/bin/python3 /home/jabberbot/urlbot/urlbot.py
directory=/home/jabberbot/urlbot/
stderr_logfile=/home/jabberbot/urlbot.err
stdout_logfile=/home/jabberbot/urlbot.log

[program:idlebot]
directory=/home/jabberbot/urlbot/
command=/home/jabberbot/botenv/bin/python3 /home/jabberbot/urlbot/idlebot.py
stderr_logfile=/home/jabberbot/idlebot.err
stdout_logfile=/home/jabberbot/idlebot.log
